# deploy/laysh.service
[Unit]
Description=Laysh Arabic-first learning simulations
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/laysh
EnvironmentFile=-%h/.config/laysh/service.env
Environment=LAYSH_HOST=127.0.0.1
Environment=LAYSH_PORT=8765
Environment=LAYSH_CODEX_BACKEND=codex
Environment=LAYSH_RECORD_RUNTIME=0
Environment=LAYSH_UNDERSTAND_MODEL=gpt-5.6-luna
Environment=LAYSH_UNDERSTAND_FALLBACK_MODEL=gpt-5.6-sol
Environment=LAYSH_EVIDENCE_UNDERSTAND_MODEL=gpt-5.6-sol
Environment=LAYSH_GENERATE_MODEL=gpt-5.6-sol
Environment=LAYSH_HEAL_MODEL=gpt-5.6-sol
Environment=LAYSH_QA_MODEL=gpt-5.6-sol
Environment=LAYSH_PUBLIC_JOB_TIMEOUT_SECONDS=180
Environment=LAYSH_PUBLIC_STAGE_TIMEOUT_SECONDS=90
Environment=LAYSH_PUBLIC_QA_TIMEOUT_SECONDS=45
ExecStart=%h/laysh/scripts/serve.sh
Restart=always
RestartSec=3s
TimeoutStopSec=20s
KillMode=control-group
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/laysh/out
ReadWritePaths=%h/.codex

[Install]
WantedBy=default.target
